# ebpf_service/ebpf.service
[Unit]
Description=Simple eBPF service
After=network.target 

[Service]
User=root
WorkingDirectory=/local/repository/ebf_service
ExecStart=/usr/bin/python3 /local/repository/ebf_service/service.py
Restart=always

[Install]
WantedBy=multi-user.target